# src/main/debian/prerm
#!/bin/bash

# Stop running instance if needed
systemctl stop gridepo.service

# Disable service if exists
systemctl disable gridepo.service

# remove symlink
rm /usr/bin/gridepo
